# supervisord.conf
[program:elasticsearch]
command=/opt/elasticsearch/bin/elasticsearch
user=elasticsearch